# backend/ansible/roles/rpi/templates/etc/supervisor/conf.d/rpi-obd.conf
[program:rpi-obd]
command=/usr/bin/python3 /root/backend/rpi/obd.py
autostart=false
autorestart=true
stderr_logfile=/var/log/rpi-obd.err.log
stdout_logfile=/var/log/rpi-obd.out.log