# backend/ansible/roles/rpi/templates/etc/supervisor/conf.d/rpi-gps.conf
[program:rpi-gps]
command=/usr/bin/python3 /root/backend/rpi/gps.py
autostart=false
autorestart=true
stderr_logfile=/var/log/rpi-gps.err.log
stdout_logfile=/var/log/rpi-gps.out.log
