# backend/ansible/roles/rpi/templates/etc/cron.hourly/dnsimple
#!/bin/bash

/usr/bin/python3 /root/backend/rpi/dnsimple_update.py
